# debian/tmp/lib/systemd/system/ddnsd-java.service
[Unit]
Description=DDNS Daemon
After=network.target

[Service]
ExecStart=/usr/bin/ddnsd-java

[Install]
WantedBy=multi-user.target\n
